# deploy/debate.service
[Unit]
Description=Debate Service
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/debate-service
EnvironmentFile=/etc/default/debate-env
ExecStart=/opt/debate-service/.venv/bin/python /opt/debate-service/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
